# deploy/fifth-grape-prune-events.service
[Unit]
Description=Fifth Grape — prune analytics events older than retention window
After=network.target

[Service]
Type=oneshot
User=fifth-grape
WorkingDirectory=/home/fifth-grape/backend
EnvironmentFile=/home/fifth-grape/backend/.env
ExecStart=/home/fifth-grape/backend/.venv/bin/python scripts/prune_events.py
